# systemd/tuxedo_keyboard.path
[Unit]
Description=Tastaturbeleuchtung ändern

[Path]
PathChanged=/home/te/tuxedo_keyboard.ini

[Install]
WantedBy=multi-user.target
